# image-cleaner.timer
[Unit]
Description=Runs docker garbage collector

[Timer]
OnBootSec=2h
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target

[X-Fleet]
Global=true

