# mongo-backup.timer
[Unit]
Description=Runs the mongo backup job

[Timer]
OnCalendar=11:00

[Install]
WantedBy=timers.target

[X-Fleet]
MachineOf=mongo-backup.service

