# wordpress-article-transformer-sidekick@.service
[Unit]
Description=wordpress article transformer sidekick
BindsTo=wordpress-article-transformer@%i.service
After=wordpress-article-transformer@%i.service

[Service]
ExecStart=/bin/sh -c "\
  etcdctl set   /ft/services/wordpress-article-transformer/healthcheck     true; \
    etcdctl mkdir /ft/services/wordpress-article-transformer/servers; \
    etcdctl set /ft/healthcheck/wordpress-article-transformer-%i/path /__health; \
  while true; do \
      PORT=`[ $PORT ] && echo $PORT || echo $(/usr/bin/docker port wordpress-article-transformer-%i 8080 | cut -d':' -f2)`; \
      etcdctl set /ft/services/wordpress-article-transformer/servers/%i "http://$HOSTNAME:$PORT" --ttl 600; \
      sleep 120; \
    done"
ExecStop=/usr/bin/etcdctl rm /ft/services/wordpress-article-transformer/servers/%i

[X-Fleet]
MachineOf=wordpress-article-transformer@%i.service
