# service-files/image-cleaner.service
[Unit]
Description=Job to clean up docker images + containers on hosts
Wants=docker.socket

[Service]
Type=oneshot
ExecStart=-/bin/sh -c "docker rm `docker ps -a -q --no-trunc`;"
ExecStart=-/bin/sh -c "docker rmi `docker images -f \"dangling=true\" -q --no-trunc`;"
ExecStart=-/bin/sh -c "rm -rf /var/lib/docker/tmp/*"

[Install]
WantedBy=multi-user.target

[X-Fleet]
Global=true

